# deploy-examples/gunicorn/webtrace.socket
[Unit]
Description=webtrace gunicorn socket

[Socket]
ListenStream=/run/webtrace.sock
SocketUser=www-data

[Install]
WantedBy=sockets.target
